# ps3remote.service
[Unit]
Description=PS3 Remote Daemon

[Service]
ExecStart=/usr/bin/ps3rd.sh

[Install]
WantedBy=multi-user.target